# data/supervisord.conf
[supervisord]
nodaemon=true

[program:candidcloud]
command=/data/entrypoint.sh
